# rsyncbackup/backup.cron
# Run system wide backup daily at 3am by default
0 3 * * * root /sbin/backup
