# gentoo/emerge-sync.cron
# Run system wide emerge and sync daily at 3am by default
0 3 * * * root /usr/local/sbin/emerge-sync
